# systemd/patchwave-report.service
[Unit]
Description=PatchWave daily Cowork report
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/patchwave/bin/patchwave-report
User=root
StandardOutput=journal
StandardError=journal
